# docker/supervisord.conf
[program:dockerd]
command=/usr/local/bin/docker-entrypoint.sh dockerd

